# deploy/intern-checkin-bot.service
[Unit]
Description=Intern Check-In Telegram Bot + Scheduler
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/intern-checkin-bot
EnvironmentFile=/opt/intern-checkin-bot/.env
ExecStart=/opt/intern-checkin-bot/.venv/bin/python main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=intern-checkin-bot

[Install]
WantedBy=multi-user.target
